# deploy/firmin-comparison.timer
[Unit]
Description=Run Firmin comparison report daily at 8am UK time

[Timer]
OnCalendar=*-*-* 08:00:00
TimeZone=Europe/London
Persistent=true

[Install]
WantedBy=timers.target
