# deploy/firmin-verification.timer
[Unit]
Description=Run Firmin verification backfill twice daily

[Timer]
OnCalendar=*-*-* 14:00:00
OnCalendar=*-*-* 22:00:00
TimeZone=Europe/London
Persistent=true

[Install]
WantedBy=timers.target
